# roles/rag/templates/chromadb.service.j2
[Unit]
Description=ChromaDB Vector Database
After=network.target

[Service]
Type=simple
User={{ rag_user }}
ExecStart=chroma run --host {{ chromadb_host }} --port {{ chromadb_port }} --path {{ chromadb_data_path }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
